# tests/supervisor-app.ini
[eventlistener:supervisor-event-listener]
command=/usr/local/bin/supervisor-event-listener 
    -c /etc/supervisor-event-listener.toml
events=
    PROCESS_STATE_EXITED,
    PROCESS_STATE_FATAL,
    PROCESS_STATE_STOPPED,
    PROCESS_STATE_RUNNING
user=nobody
group=nobody

[program:sleep-then-exit]
command=sleep 3 && exit 1
